# etc/systemd/system/dynv6.timer
[Unit]
Description=Run dynv6 hourly

[Timer]
OnCalendar=hourly
Unit=dynv6.service

[Install]
WantedBy=timers.target
